# .config/autostart/mjpg_streamer.desktop
[Desktop Entry]
Name=mjpg_streamer
Exec=lxterminal -e "/home/pi/Desktop/run_mjpg_streamer.sh"
Type=Application
