# .config/autostart/nodejs_server.desktop
[Desktop Entry]
Name=nodejs_server
Exec=lxterminal -e "/home/pi/Desktop/run_nodejs_srv.sh"
Type=Application

